# tools/miniconda3/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/timwood/Desktop/projects/PWA/MV/tools/miniconda3/bin/conda";
setenv _CONDA_EXE "/Users/timwood/Desktop/projects/PWA/MV/tools/miniconda3/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/timwood/Desktop/projects/PWA/MV/tools/miniconda3/bin/python";
setenv _CONDA_ROOT "/Users/timwood/Desktop/projects/PWA/MV/tools/miniconda3";

source "/Users/timwood/Desktop/projects/PWA/MV/tools/miniconda3/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"